# printer_capture.service
[Unit]
Description=USB bridge printer capture
After=multi-user.target

[Service]
Type=simple
Environment=PRINTER_DEV=/dev/g_printer0
Environment=PRINT_IDLE_TIMEOUT=2.0
ExecStart=/usr/bin/python3 /root/usb_bridge/bin/printer_capture.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
