# root/usb_bridge/systemd/report_uploader.service
[Unit]
Description=USB bridge report converter and uploader
After=network-online.target printer_capture.service
Wants=network-online.target

[Service]
Type=simple
Environment=GS_BIN=gs
Environment=REPORT_UPLOAD_URL=http://8.148.73.190:5000/upload
ExecStart=/usr/bin/python3 /root/usb_bridge/bin/report_uploader.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
