# root/usb_bridge/systemd/scan_patient_capture.service
[Unit]
Description=USB bridge scanner patient capture
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /root/usb_bridge/bin/scan_patient_capture.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
